# build/files/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock   ; (the path to the socket file)

[supervisord]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
loglevel=debug   

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; ip_address:port specifier, *:port for all iface

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock 

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:zerotier]
command=zerotier-one /zt-data
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:dcloudmonitor]
command=node /usr/dcloudmonitor/index.js "/data/acloud/config" /etc/haproxy/newconfig /zt-data/authtoken.secret_ 
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true